# config/supervisor.conf.in
[group:muchmala]
programs=frontend, io, app, proxy, stylus


[program:frontend]
process_name=%(program_name)s_807%(process_num)s
numprocs=<%= config.serversCount.frontend %>
command=node index.js --port 807%(process_num)s
directory=<%= root %>/components/muchmala-frontend

autostart=<%= config.autoRestart %>
autorestart=<%= config.autoRestart %>


[program:io]
process_name=%(program_name)s_809%(process_num)1d
numprocs=<%= config.serversCount.io %>
command=node index.js --externalHost io%(process_num)s.<%= config.http.hostName %> --internalPort 809%(process_num)s
directory=<%= root %>/components/muchmala-io

autostart=<%= config.autoRestart %>
autorestart=<%= config.autoRestart %>


[program:app]
process_name=%(program_name)s_%(process_num)02d
numprocs=<%= config.serversCount.app %>
command=node index.js
directory=<%= root %>/components/muchmala-app

autostart=<%= config.autoRestart %>
autorestart=<%= config.autoRestart %>


[program:proxy]
command=<%= root %>/bin/proxy -c <%= root %>/proxy.json

autostart=true
autorestart=unexpected


[program:stylus]
command=jake stylus-watch
directory=<%= root %>/components/muchmala-frontend

autostart=true
autorestart=unexpected